# lore/extensions/cc-telemetry/launchd/com.layerdynamics.cc-telemetry.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.layerdynamics.cc-telemetry</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__CC_TELEMETRY_DAEMON_DIR__/daemon.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.claude/telemetry/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/telemetry/daemon.log</string>

    <key>WorkingDirectory</key>
    <string>__CC_TELEMETRY_DAEMON_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>CC_TELEMETRY_LOG_LEVEL</key>
        <string>INFO</string>
    </dict>

    <!-- Throttle restarts to avoid CPU spin on repeated crashes -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
